# stable-audio-sfx-mcp.service
[Unit]
Description=Stable Audio 3 SFX Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/projects/stable-audio-3-sfx
ExecStart=%h/projects/stable-audio-3-sfx/.venv/bin/python %h/projects/stable-audio-3-sfx/main.py
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=120
StartLimitBurst=3
Environment=PYTHONUNBUFFERED=1
Environment=PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True

[Install]
WantedBy=default.target
